# docker_setup/deploy/eta.service
[Unit]
Description=ETA infrastructure service
Requires=docker.service
After=docker.service

[Service]
Restart=on-failure
WorkingDirectory=/opt/intel/eta/
ExecStart=/opt/intel/eta/compose_startup.sh --up
ExecStop=/opt/intel/eta/compose_startup.sh --down
TimeoutSec=301s

[Install]
WantedBy=multi-user.target
